# scripts/launchd/com.mitchell.stack-ops.cheap-batch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  cheap-batch: route recurring bulk toil to the open-weight tier, unattended.

  SHIPS DOUBLE-DISABLED, matching the career-ops convention for every unattended
  job that spends money:
    1. Disabled=true below, and
    2. the CHEAP_BATCH_ENABLED kill switch, which this plist deliberately does
       NOT set. The --scheduled flag makes the script a 0-exit no-op without it.
  Both gates must be lifted by hand. Verified 2026-07-22: with the switch unset,
  a --scheduled --apply run exits 0 having called nothing.

  Arm only after a live verification:
    node scripts/cheap-batch.mjs                 # DRY, calls nothing
    node scripts/cheap-batch.mjs --apply         # one real run, check the output
    printf '\nexport CHEAP_BATCH_ENABLED=true\n' >> ~/.zshenv
    cp scripts/launchd/com.mitchell.stack-ops.cheap-batch.plist ~/Library/LaunchAgents/
    /usr/libexec/PlistBuddy -c "Set :Disabled false" ~/Library/LaunchAgents/com.mitchell.stack-ops.cheap-batch.plist
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.mitchell.stack-ops.cheap-batch.plist
    launchctl enable gui/$UID/com.mitchell.stack-ops.cheap-batch

  THE ENV VAR MUST GO IN THE LOGIN SHELL (~/.zshenv), NOT A PROJECT .env.
  launchd jobs inherit no shell environment, and this script reads the switch
  from process.env. A project .env would leave the job no-opping FOREVER with no
  error, because a disabled kill switch exits 0 by design. The only reason a
  login-shell export works at all is that ProgramArguments invokes `zsh -lc`.
  Verified 2026-07-22 by exporting it in ~/.zshenv and running the exact
  ProgramArguments command: the run got past the kill switch.

  Also create the log directory first, or launchd fails SILENTLY on the
  StandardOutPath/StandardErrorPath below:
    mkdir -p ~/Library/Logs/stack-ops

  Note: `launchctl disable` persists across bootstraps. If a bootstrap returns
  "Input/output error" with a valid plist, check `launchctl print-disabled`.

  Runs 04:00 PT daily, deliberately after career-ops' 03:00 Sunday drain so the
  two never contend for the same logs.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mitchell.stack-ops.cheap-batch</string>

  <key>Disabled</key>
  <true/>

  <!--
    Invoked through a login shell on purpose. node lives under nvm here
    (~/.nvm/versions/node/<ver>/bin/node), not /usr/local/bin/node, and launchd
    jobs get a minimal PATH that contains neither. Hardcoding the nvm path would
    work today and break silently on the next node upgrade, at 04:00, with the
    only evidence in a log nobody reads. `zsh -lc` resolves node the same way an
    interactive shell does, so a version bump is a non-event.
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec node /Users/mitchellwilliams/Documents/stack-ops/scripts/cheap-batch.mjs --apply --scheduled</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mitchellwilliams/Documents/stack-ops</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/mitchellwilliams/Library/Logs/stack-ops/cheap-batch.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mitchellwilliams/Library/Logs/stack-ops/cheap-batch.err</string>

  <key>RunAtLoad</key>
  <false/>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
